# conf/crontab.cron
*/30 * * * * /Users/pr/work/code/py/wsd/wr.sh
