# conf/crontab.cron
*/30 * * * * /User/pr/work/code/py/wsd/wr.sh
